# dotfiles/.config/systemd/user/reset-qutebrowser.service
[Unit]
Description=Reset qutebrowser data

[Service]
ExecStart=%h/.local/bin/reset-qutebrowser

[Install]
WantedBy=default.target
