# dotfiles/.config/systemd/user/reset-qutebrowser.service
[Unit]
Description=Reset qutebrowser data

[Service]
ExecStart=reset-qutebrowser

[Install]
WantedBy=default.target
